# ops/com.media-studio.web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.media-studio.web</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>__MEDIA_STUDIO_ROOT__/scripts/start_web.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__MEDIA_STUDIO_ROOT__</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>NODE_ENV</key>
      <string>production</string>
      <key>PORT</key>
      <string>3000</string>
      <key>MEDIA_STUDIO_SUPERVISOR</key>
      <string>launchd</string>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/media-studio-web.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/media-studio-web.error.log</string>
  </dict>
</plist>
